# systemd-units/giessknecht.service
[Unit]
Description=Execute irrigation cycle

[Service]
Type=oneshot
ExecStart=/usr/local/bin/giessknecht.py schedule.json
